# scripts/systemd/edgehunter-worker.service
[Unit]
Description=EdgeHunter Async Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/EdgeHunter
EnvironmentFile=%h/EdgeHunter/.env
Environment=PYTHONPATH=%h/EdgeHunter
ExecStart=%h/EdgeHunter/.venv/bin/python %h/EdgeHunter/backend/async_worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
